# example/crontabs/root
* * * * * curl -s4 ifconfig.co > /dev/stdout
